# playbooks/roles/promtail/templates/promtail.service.j2
{% if container_command == 'podman' %}
# promtail.service
# autogenerated by Podman 3.2.3
# Mon Sep 27 07:11:04 UTC 2021

[Unit]
Description=Podman promtail.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70

ExecStartPre=/bin/rm -f %t/container-promtail.pid %t/container-promtail.ctr-id

ExecStart=/usr/bin/podman run \
          --conmon-pidfile %t/container-promtail.pid \
          --cidfile %t/container-promtail.ctr-id \
          --cgroups=no-conmon \
          -d \
          --name promtail \
          --log-opt=path=/dev/null \
          -v /etc/promtail/config.yml:/etc/promtail/config.yml:ro \
          -v /var/log/journal:/var/log/journal:ro,z \
          -v /var/log/audit:/var/log/audit:ro,z \
          -v promtail-run:/run/promtail:rw,z \
          -p 3101:3101 \
          {{ promtail_image }}

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-promtail.ctr-id -t 10

ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-promtail.ctr-id

PIDFile=%t/container-promtail.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target

{% else %}
# promtail.service

[Unit]
Description=promtail.service
Wants=network.target
After=network-online.target

[Service]
Restart=on-failure
TimeoutStopSec=70

ExecStart=/usr/bin/docker run  \
          --name promtail \
          -v /etc/promtail/config.yml:/etc/promtail/config.yml:ro \
          -v /var/log/journal:/var/log/journal:ro,z \
          -v /var/log/audit:/var/log/audit:ro,z \
          -v promtail-run:/run/promtail:rw,z \
          -p 3101:3101 \
          {{ promtail_image }}

ExecStop=/usr/bin/docker stop promtail -t 10

ExecStopPost=/usr/bin/docker rm -f promtail

[Install]
WantedBy=multi-user.target default.target
{% endif %}
