# playbooks/roles/carbonapi/templates/carbonapi-healthcheck.timer.j2
[Unit]
Description=Runs graphite healthcheck every 5 minutes
Requires=graphite-healthcheck.service

[Timer]
Unit=graphite-healthcheck.service
OnUnitInactiveSec=5m
AccuracySec=1s

[Install]
WantedBy=timers.target

